# docker/config/supervisor/conf.d/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
; command=apachectl -D "FOREGROUND" -k start
username=www-data
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
